# contentmanagementservice/conf/supervisor.conf
[supervisord]
childlogdir = %(here)s/../../logs/
logfile = %(here)s/../../logs/supervisord.log
pidfile = /tmp/supervisord-cms.pid
logfile_maxbytes = 50MB
logfile_backups=100
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
# user = 1000
identifier = supervisor
nocleanup = true
strip_ansi = false
# environment = KEY1="value1",KEY2="value2"

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = 8481

[supervisorctl]
serverurl = http://localhost:8481

[program:contentmanagementservice]
directory=%(here)s/../conf
command=gunicorn -c gunicorn.py service_app:app
autorestart=true
stdout_logfile=%(here)s/../../logs/cms-stdout.log
stderr_logfile=%(here)s/../../logs/cms-stderr.log